# context/automation/launchd/ai.aoi.ssot.monthly_security_audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
<key>Label</key><string>ai.aoi.ssot.monthly_security_audit</string>
<key>ProgramArguments</key><array>
<string>/Users/silkroadcat/.openclaw/workspace/scripts/_tg_alert_on_fail.sh</string>
<string>monthly_security_audit</string>
<string>/Users/silkroadcat/.openclaw/workspace/scripts/monthly_security_audit.sh</string>
</array>
<key>StartCalendarInterval</key><dict><key>Day</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>11</integer></dict>
<key>WorkingDirectory</key><string>/Users/silkroadcat/.openclaw/workspace</string>
<key>StandardOutPath</key><string>/tmp/launchd.ai.aoi.ssot.monthly_security_audit.out.log</string>
<key>StandardErrorPath</key><string>/tmp/launchd.ai.aoi.ssot.monthly_security_audit.err.log</string>
</dict></plist>
